# build/librenms/files/etc/service/memcached/run
#!/bin/bash

echo "STARTING MEMCACHED"

/usr/bin/memcached -m 64 -u root
